# meta-citadel/recipes-citadel/citadel-config/files/systemd/watch-run-user.path
[Unit]
Description=Watch for creation of /run/user/1000
After=run-user-1000.mount
Requires=run-user-1000.mount

[Path]
PathExists=/run/user/1000

[Install]
WantedBy=multi-user.target
