# meta-citadel/recipes-core/launch-default-realm/files/watch-run-user.service
[Unit]
Description=Watch run-user service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --no-block start launch-default-realm.path
