# meta-gnome/recipes-connectivity/networkmanager/files/watch-resolvconf.path
[Unit]
Description=Watch /run/NetworkManager

[Path]
PathExists=/run/NetworkManager/resolv.conf
PathChanged=/run/NetworkManager
MakeDirectory=true

[Install]
WantedBy=multi-user.target
