# meta-citadel/recipes-initrd/citadel-initramfs/files/citadel-rootfs-mount.path
[Unit]
Description=Watch for new entries in /dev/mapper

[Path]
PathChanged=/dev/mapper

