# meta-citadel/recipes-initrd/citadel-initramfs/files/citadel-rootfs-mount.service
[Unit]
Description=Mount citadel rootfs

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/sh -c 'echo "hello world"'
