# meta-citadel/recipes-initrd/citadel-initramfs/files/citadel-rootfs-mount.service
[Unit]
Description=Mount citadel rootfs

ConditionPathExists=/dev/mapper/citadel-rootfsA
ConditionPathExists=/dev/mapper/citadel-rootfsB
ConditionPathExists=/dev/mapper/citadel-storage

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c 'echo "hello world"'
