# meta-citadel/recipes-citadel/launch-default-realm/files/launch-default-realm.path
[Unit]
Description=Wait for creation of /run/user/1000/pulse/native

[Path]
PathExists=/run/user/1000/pulse/native
