# meta-citadel/recipes-core/primary-user-appimg/files/watch-run-user.service
[Unit]
Description=Watch run-user service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --no-block start primary-user-appimg.path
